# .config/profile.d/beets.sh
# Add projects to python path
export PYTHONPATH="${HOME}/Projects/sbp-beets:${PYTHONPATH}"
