# .config/profile.d/lf.sh
# Source our function
LFCD="${HOME}/.config/lf/lfcd.sh"
if [ -f "$LFCD" ]; then
   source "$LFCD"
fi
