# .config/profile.d/r.sh
# Rstudio stuff
export R_PROFILE_USER="${XDG_CONFIG_HOME}/R/Rprofile"
export R_ENVIRON_USER="${XDG_CONFIG_HOME}/R/Renviron"
export R_LIBS_USER="${HOME}/Data/R/%p-library/%v"
